# infra/cron/riga-guide.cron
# Riga Guide — automated tasks
55 23 * * * cd /root/Riga_guide && /usr/bin/docker compose run --rm ingest python scripts/daily_rollup.py --out logs/rollup/$(date +\%Y-\%m-\%d).md >> logs/cron.log 2>&1
0 4 * * * cd /root/Riga_guide && /bin/bash scripts/backup_kb.sh >> logs/cron.log 2>&1
